# server/tm2.service
[Unit]
Description=TM2 FastAPI
After=network.target

[Service]
User=floindil
Group=floindil
WorkingDirectory=/home/floindil/vapp/TableMaker2
Environment="PATH=/home/floindil/vapp/TableMaker2/.venv/bin"
ExecStart=/home/floindil/vapp/TableMaker2/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target